# autoStart_node.service
[Unit]
Description = starting node-js
After=autoStart_python.service

[Service]
ExecStart=/home/pi/bezelie/dev_edgar/boot_node.sh
Type=simple
Restart=no
User=pi
Group=pi

[Install]
WantedBy=multi-user.target
